# .config/systemd/user/paccache-local.service
[Unit]
Description=Remove unused cached package files

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -c ${HOME}/.cache/pacman/pkg -rk1
