# .config/systemd/user/mpd.service.d/override.conf
[Unit]
Upholds=mpd-notification.service
