# .config/systemd/user/mpd-notification.service.d/override.conf
[Unit]
 After=dunst.service xfce4-notifyd.service mako.service

[Service]
ExecStart=
ExecStart=/usr/bin/mpd-notification  -H localhost -p 6600 -m '${HOME}/Music' c -s 150 -t 10
